# docker/config/supervisord-headful.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

; Control socket so `supervisorctl status` / `restart <prog>` work inside the
; container — without it a single dead service can't be restarted or inspected
; without nuking the whole container.
[unix_http_server]
file=/run/supervisor.sock

[rpcinterface:supervisor]
; Entry-point form (module:attr) — Debian's supervisor resolves the factory
; via importlib.metadata, which rejects the classic dotted module.attr spec.
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

; System dbus — without it Chrome stalls on dbus autolaunch attempts.
[program:dbus]
command=dbus-daemon --system --nofork --nopidfile --address=unix:path=/run/dbus/system_bus_socket
autostart=true
autorestart=true
stdout_logfile=/var/log/dbus.log
stderr_logfile=/var/log/dbus-err.log
priority=5

[program:sway]
command=sway %(ENV_SWAY_EXTRA_ARGS)s
environment=
    XDG_RUNTIME_DIR="/tmp/xdg-runtime",
    WLR_BACKENDS="headless",
    WLR_HEADLESS_OUTPUTS="2",
    WLR_LIBINPUT_NO_DEVICES="1",
    WLR_NO_HARDWARE_CURSORS="1",
    WLR_RENDERER="%(ENV_WLR_RENDERER)s",
    LIBSEAT_BACKEND="noop"
autostart=true
autorestart=true
stdout_logfile=/var/log/sway.log
stderr_logfile=/var/log/sway-err.log
priority=10

; One wayvnc per output. Capture capped at 30 fps so VNC encoding
; doesn't contend with Chrome on the compositor's render loop.
[program:wayvnc]
command=wayvnc --gpu --max-fps=30 --output=HEADLESS-1 -S /tmp/wayvnc-1.sock 127.0.0.1 %(ENV_VNC_PORT)s
environment=
    XDG_RUNTIME_DIR="/tmp/xdg-runtime",
    WAYLAND_DISPLAY="wayland-1"
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/wayvnc.log
stderr_logfile=/var/log/wayvnc-err.log
priority=20

[program:wayvnc-2]
command=wayvnc --gpu --max-fps=30 --output=HEADLESS-2 -S /tmp/wayvnc-2.sock 127.0.0.1 %(ENV_VNC_PORT_2)s
environment=
    XDG_RUNTIME_DIR="/tmp/xdg-runtime",
    WAYLAND_DISPLAY="wayland-1"
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/wayvnc-2.log
stderr_logfile=/var/log/wayvnc-2-err.log
priority=20

[program:novnc]
command=websockify --web=/usr/share/novnc 127.0.0.1:%(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc-err.log
priority=25

[program:novnc-2]
command=websockify --web=/usr/share/novnc 127.0.0.1:%(ENV_NOVNC_PORT_2)s localhost:%(ENV_VNC_PORT_2)s
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/novnc-2.log
stderr_logfile=/var/log/novnc-2-err.log
priority=25

; Headful Chrome is the human-parity path. Keep this launch surface close to
; nodriver/ordinary Chrome: compositor/GPU flags needed for Wayland in Docker,
; profile/first-run hygiene, and a fixed CDP port only. Avoid broad background,
; component, renderer-throttling, IPC, and blink-feature suppression flags here;
; those improve crawler throughput but make the browser less human-shaped.
[program:chrome-1]
command=google-chrome-stable
    --no-sandbox
    --class=chrome-1
    --ozone-platform=wayland
    --enable-features=UseOzonePlatform
    --ignore-gpu-blocklist
    --enable-gpu-rasterization
    --enable-zero-copy
    --disable-dev-shm-usage
    --disable-infobars
    --disable-session-crashed-bubble
    --disable-search-engine-choice-screen
    --no-first-run
    --no-service-autorun
    --no-default-browser-check
    --no-pings
    --password-store=basic
    --homepage=about:blank
    --remote-debugging-port=%(ENV_CDP_PORT_1)s
    --user-data-dir=%(ENV_CHROME_PROFILE_DIR_1)s
environment=
    XDG_RUNTIME_DIR="/tmp/xdg-runtime",
    WAYLAND_DISPLAY="wayland-1",
    DBUS_SESSION_BUS_ADDRESS="unix:path=/run/dbus/system_bus_socket",
    DISPLAY=""
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/chrome-1.log
stderr_logfile=/var/log/chrome-1-err.log
priority=30

[program:chrome-2]
command=google-chrome-stable
    --no-sandbox
    --class=chrome-2
    --ozone-platform=wayland
    --enable-features=UseOzonePlatform
    --ignore-gpu-blocklist
    --enable-gpu-rasterization
    --enable-zero-copy
    --disable-dev-shm-usage
    --disable-infobars
    --disable-session-crashed-bubble
    --disable-search-engine-choice-screen
    --no-first-run
    --no-service-autorun
    --no-default-browser-check
    --no-pings
    --password-store=basic
    --homepage=about:blank
    --remote-debugging-port=%(ENV_CDP_PORT_2)s
    --user-data-dir=%(ENV_CHROME_PROFILE_DIR_2)s
environment=
    XDG_RUNTIME_DIR="/tmp/xdg-runtime",
    WAYLAND_DISPLAY="wayland-1",
    DBUS_SESSION_BUS_ADDRESS="unix:path=/run/dbus/system_bus_socket",
    DISPLAY=""
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/chrome-2.log
stderr_logfile=/var/log/chrome-2-err.log
priority=30
